# puppet/modules/supervisor/files/supervisor.conf
description     "supervisor"

setuid vagrant
setgid supervisor

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/local/bin/supervisord -c /etc/supervisord.conf
